# roles/servarr_stack/templates/qbittorrent2.service.j2
[Unit]
Description=qBittorrent-nox instance for qbittorrent2
After=network.target

[Service]
User=qbittorrent1
Group=media
ExecStart=/usr/bin/qbittorrent-nox --webui-port=8081
Restart=on-failure

[Install]
WantedBy=multi-user.target
