# roles/servarr_stack/templates/qbittorrent.service.j2
[Unit]
Description=qBittorrent-nox service
After=network.target

[Service]
User=qbittorrent
Group=media
ExecStart=/usr/bin/qbittorrent-nox
Restart=on-failure

[Install]
WantedBy=multi-user.target
